# docker/services.d/postgres/run
#!/command/with-contenv sh
exec s6-setuidgid postgres /usr/lib/postgresql/17/bin/postgres -D /var/lib/postgresql/17/main
